# obexd/.footprint
drwxr-xr-x	timcowchip/users	etc/
drwxr-xr-x	timcowchip/users	etc/rc.d/
-rwxr-xr-x	timcowchip/users	etc/rc.d/obexd
drwxr-xr-x	timcowchip/users	usr/
drwxr-xr-x	timcowchip/users	usr/libexec/
drwxr-xr-x	timcowchip/users	usr/libexec/obex/
-rwxr-xr-x	timcowchip/users	usr/libexec/obex/obex-client
-rwxr-xr-x	timcowchip/users	usr/libexec/obex/obexd
drwxr-xr-x	timcowchip/users	usr/share/
drwxr-xr-x	timcowchip/users	usr/share/dbus-1/
drwxr-xr-x	timcowchip/users	usr/share/dbus-1/services/
-rw-r--r--	timcowchip/users	usr/share/dbus-1/services/obex-client.service
-rw-r--r--	timcowchip/users	usr/share/dbus-1/services/obexd.service
